# contrib/init/quantumd.conf
description "Quantum Core Daemon"

start on runlevel [2345]
stop on starting rc RUNLEVEL=[016]

env QUANTUMD_BIN="/usr/bin/quantumd"
env QUANTUMD_USER="quantum"
env QUANTUMD_GROUP="quantum"
env QUANTUMD_PIDDIR="/var/run/quantumd"
# upstart can't handle variables constructed with other variables
env QUANTUMD_PIDFILE="/var/run/quantumd/quantumd.pid"
env QUANTUMD_CONFIGFILE="/etc/quantum/quantum.conf"
env QUANTUMD_DATADIR="/var/lib/quantumd"

expect fork

respawn
respawn limit 5 120
kill timeout 60

pre-start script
    # this will catch non-existent config files
    # quantumd will check and exit with this very warning, but it can do so
    # long after forking, leaving upstart to think everything started fine.
    # since this is a commonly encountered case on install, just check and
    # warn here.
    if ! grep -qs '^rpcpassword=' "$QUANTUMD_CONFIGFILE" ; then
        echo "ERROR: You must set a secure rpcpassword to run quantumd."
        echo "The setting must appear in $QUANTUMD_CONFIGFILE"
        echo
        echo "This password is security critical to securing wallets "
        echo "and must not be the same as the rpcuser setting."
        echo "You can generate a suitable random password using the following"
        echo "command from the shell:"
        echo
        echo "bash -c 'tr -dc a-zA-Z0-9 < /dev/urandom | head -c32 && echo'"
        echo
        echo "It is also recommended that you also set alertnotify so you are "
        echo "notified of problems:"
        echo
        echo "ie: alertnotify=echo %%s | mail -s \"Quantum Alert\"" \
            "admin@foo.com"
        echo
        exit 1
    fi

    mkdir -p "$QUANTUMD_PIDDIR"
    chmod 0755 "$QUANTUMD_PIDDIR"
    chown $QUANTUMD_USER:$QUANTUMD_GROUP "$QUANTUMD_PIDDIR"
    chown $QUANTUMD_USER:$QUANTUMD_GROUP "$QUANTUMD_CONFIGFILE"
    chmod 0660 "$QUANTUMD_CONFIGFILE"
end script

exec start-stop-daemon \
    --start \
    --pidfile "$QUANTUMD_PIDFILE" \
    --chuid $QUANTUMD_USER:$QUANTUMD_GROUP \
    --exec "$QUANTUMD_BIN" \
    -- \
    -pid="$QUANTUMD_PIDFILE" \
    -conf="$QUANTUMD_CONFIGFILE" \
    -datadir="$QUANTUMD_DATADIR" \
    -disablewallet \
    -daemon

